# packages/virt-firmware/virt_firmware-25.12.tar.gz/virt_firmware-25.12/systemd/kernel-bootcfg-boot-successful.service
#
# SPDX-License-Identifier: LGPL-2.1-or-later
# (c) 2023 Vitaly Kuznetsov
#

[Unit]
Description=UKI Successful Boot
Requires=boot-complete.target
After=local-fs.target boot-complete.target
Conflicts=shutdown.target
Before=shutdown.target
ConditionFirmware=uefi
ConditionPathExists=!/etc/initrd-release

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=kernel-bootcfg --boot-successful --update-csv

[Install]
WantedBy=basic.target
