# packages/plex-letterboxd/plex_letterboxd-0.4.tar.gz/plex_letterboxd-0.4/systemd/plex-letterboxd@.timer
[Unit]
Description=Plex to Letterboxd export timer (%i)

[Timer]
OnCalendar=%i
Persistent=true
RandomizedDelaySec=30m
Unit=plex-letterboxd@%i.service

[Install]
WantedBy=timers.target
