# packages/fujin-cli/fujin_cli-0.14.1.tar.gz/fujin_cli-0.14.1/src/fujin/templates/default.timer.j2
[Unit]
Description=Timer for {{ process_name }}

[Timer]
OnCalendar={{ process.timer }}
Unit={{ process_name }}.service
Persistent=true

[Install]
WantedBy=timers.target
