# packages/tess-publisher/tess_publisher-0.9.1.tar.gz/tess_publisher-0.9.1/files/tesspublisher.service
[Unit]
Description=TESS Quality Sky Meter MQTT publisher service
Documentation=https://github.com/STARS4ALL/tess-publisher

[Service]
Type=simple
User=root
KillMode=process
ExecStart=/home/pi/tess/.venv/bin/tess-publisher --config /home/pi/tess/config.toml --log-file /home/pi/tess/log/tess-publisher.log
EnvironmentFile=/home/pi/tess/.env
WorkingDirectory=/home/pi/tess/

[Install]
WantedBy=multi-user.target
