# packages/procustodibus_agent/procustodibus_agent-1.10.2.tar.gz/procustodibus_agent-1.10.2/docker/agent-fs/etc/init.d/procustodibus-agent
#!/sbin/openrc-run

command=/usr/local/bin/procustodibus-agent
command_args="--loop=120 --verbosity=INFO"
command_background=true
description="Pro Custodibus Agent"
error_log=/dev/stderr
output_log=/dev/stdout
pidfile=/var/run/procustodibus-agent.pid

depend() {
    need localmount
    need net
    after wg-quick
}
