# packages/procustodibus_agent/procustodibus_agent-1.10.2.tar.gz/procustodibus_agent-1.10.2/etc/launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>procustodibus-agent</string>
        <key>ProgramArguments</key>
        <array>
            <string>/opt/venvs/procustodibus-agent/bin/procustodibus-agent</string>
            <string>--loop</string>
            <string>120</string>
            <string>--verbosity</string>
            <string>INFO</string>
        </array>
        <key>KeepAlive</key>
        <dict>
            <key>NetworkState</key>
            <true/>
        </dict>
        <key>EnvironmentVariables</key>
        <dict>
            <key>PATH</key>
            <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        </dict>
        <key>StandardErrorPath</key>
        <string>/var/log/procustodibus-agent.err</string>
        <key>StandardOutPath</key>
        <string>/var/log/procustodibus-agent.log</string>
    </dict>
</plist>
