# packages/qr-payment-cz/qr_payment_cz-0.2.17-py3-none-any.whl/qr_payment_cz/qr-payment-cz-server.service
[Unit]
Description=QR Payment CZ - server APP

[Service]
ExecStart=/bin/bash -c \"/usr/local/bin/qr-payment-cz-server -a "670100-2208148433/6210" --port 5001 -r "/qr"  >> /var/log/qr-payment-cz-server.log 2>&1 &\"
ExecStop=killall qr-payment-cz-server
ExecStatus=ps -ax | grep qr-payment-cz-server
Type=forking

[Install]
WantedBy=multi-user.target
