# packages/procustodibus_agent/procustodibus_agent-1.10.2.tar.gz/procustodibus_agent-1.10.2/docker/dev-fs/etc/init.d/dev-send-traffic
#!/sbin/openrc-run

command=/usr/local/bin/dev-send-traffic.sh
command_background=true
description="Dev Send-Traffic Script"
error_log=/dev/stderr
output_log=/dev/stdout
pidfile=/var/run/dev-send-traffic.pid

depend() {
    need localmount
    need net
    after procustodibus-agent
}
