# packages/procustodibus_agent/procustodibus_agent-1.10.2.tar.gz/procustodibus_agent-1.10.2/docker/dev-fs/etc/init.d/dev-down-up
#!/sbin/openrc-run

command=/usr/local/bin/dev-down-up.sh
command_background=true
description="Dev Down-Up Script"
error_log=/dev/stderr
output_log=/dev/stdout
pidfile=/var/run/dev-down-up.pid

depend() {
    need localmount
    need net
    after procustodibus-agent
}
