# packages/procustodibus_agent/procustodibus_agent-1.10.2.tar.gz/procustodibus_agent-1.10.2/docker/dev-fs/etc/init.d/dev-receive-traffic
#!/sbin/openrc-run

command=nc
command_args="-klp ${TRAFFIC_PORT:-1234}"
command_background=true
description="Dev Receive-Traffic Script"
error_log=/dev/stderr
output_log=/dev/null
pidfile=/var/run/dev-receive-traffic.pid

depend() {
    need localmount
    need net
}
