# packages/procustodibus_agent/procustodibus_agent-1.10.2.tar.gz/procustodibus_agent-1.10.2/docker/dev-fs/etc/init.d/dev-custom-hosts
#!/sbin/openrc-run

description="Dev Custom Hosts"

depend() {
    need localmount
    need net
}

start() {
    files="$(find /usr/local/share/hostnames -type f 2>/dev/null || echo '')"
    for file in $files; do
        ebegin "Applying dev custom hosts $file"
        echo -e "\n# $file" | cat - $file >>/etc/hosts
    done
    eend 0
}
